# deploy/dca-tp@.service
[Unit]
Description=DCA trailing-TP manager (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=forge
# The script reads .env from this directory (BINANCE_API_KEY / BINANCE_SECRET_KEY).
WorkingDirectory=/home/forge/scripts/trade-binance-websocket-orderbook-dca-grid
# %i is the symbol, e.g. dca-tp@SOLUSDT. Manages the trailing TP only (no grid).
ExecStart=/usr/bin/python3 -u /home/forge/scripts/trade-binance-websocket-orderbook-dca-grid/orderbook_dca_grid.py %i --tp-only
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
