# deploy/dca-futures-tp@.service
[Unit]
Description=Futures trailing-TP manager (%i) — exit only, no grid
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=forge
# The script reads .env from this directory (BINANCE_API_KEY / BINANCE_SECRET_KEY).
WorkingDirectory=/opt/trade-binance-websocket-orderbook-dca-grid
# %i is the symbol, e.g. dca-futures-tp@SOLUSDT. Manages the trailing TP only (no grid).
ExecStart=/usr/bin/python3 -u /opt/trade-binance-websocket-orderbook-dca-grid/orderbook_dca_grid.py %i --tp-only
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
